# ytextr.service
[Unit]
Description=YTaudioExtractApp
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node <path-to-install-location>/extrsvr.js
WorkingDirectory=<path-to-install-location>
Restart=on-failure

[Install]
WantedBy=multi-user.target
